# deploy/ansible/templates/netplay-server.service.j2
# Managed by Ansible (deploy/ansible/playbook.yml) — do not edit by hand.
[Unit]
Description=Reversi netplay relay server
Documentation=https://github.com/jacoberrol/reversi
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ netplay_user }}
Group={{ netplay_user }}
# Bind loopback only: the exe.dev proxy terminates TLS and forwards to this port.
ExecStart={{ netplay_bin }} 127.0.0.1:{{ netplay_port }}
EnvironmentFile={{ netplay_env_file }}
# systemd creates /var/lib/netplay (owned by the service user, writable despite
# ProtectSystem=strict) and the server keeps its SQLite database there.
StateDirectory=netplay
Environment=NETPLAY_DB={{ netplay_db }}
# Requests with this Host (via the proxy's X-Forwarded-Host) hit the REST admin API.
Environment=NETPLAY_ADMIN_HOST={{ netplay_admin_host }}
Restart=on-failure
RestartSec=2

# Hardening — the relay needs a TCP socket, its env file, and its state dir.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
